# config/supervisord.ini
[program:webhook_nodejs]
environment=NODE_ENV=production,PRODUCTION=1
command=/home/designcc/.virtualenvs/webhook/bin/node server.js
directory=/home/designcc/sites/webhook
user=designcc
autostart=true
autorestart=true
redirect_stderr=true

[program:webhook_event_saver]
command=/home/designcc/.virtualenvs/webhook/bin/python marketing_event_saver.py
directory=/home/designcc/sites/webhook/test/python
user=designcc
autostart=true
autorestart=true
redirect_stderr=true
